# .config/systemd/user/power-monitor.service
[Service]
Environment=STARTUP_WAIT=10s
ExecStart=/usr/bin/bash -c 'exec $HOME/.local/bin/power-monitor.sh'

[Install]
WantedBy=default.target
